# systemd/data-collector-9api.service
[Unit]
Description=ULTRATHINK 9-API Data Collector
After=network.target redis.service
Wants=redis.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/cashmachine
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=/opt/cashmachine/.env
ExecStart=/usr/bin/python3 /opt/cashmachine/ultrathink_data_collector_9apis.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Resource limits
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target