# systemd/data-collector.service
[Unit]
Description=CASHMACHINE Data Collector
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/cashmachine
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
ExecStart=/usr/bin/python3 /opt/cashmachine/data-collector/data_collector_secure.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target