# systemd/ultrathink.service
[Unit]
Description=ULTRATHINK Learning Engine
After=network.target redis.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/cashmachine
ExecStart=/usr/bin/python3 /home/ubuntu/ultrathink_learning_engine.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target